# docker/cron/crontab
0 0 * * * php /var/www/scripts/update.php 
0 6 * * * php /var/www/scripts/unverifiedAccountsCleanup.php