# docker/cron/crontab
* * * * * php /var/www/scripts/update.php 